# aws/deploy/config/cronlog/IntelAnlytics.crontab
0 3 * * 0 hadoop /usr/sbin/logrotate -f /etc/logrotate.d/IntelAnalytics

